# .config/systemd/user/rgb.service
[Unit]
Description=RGB
After=default.target

[Service]
ExecStart=/home/ben/.bin/set-case-rgb
Restart=on-failure

[Install]
WantedBy=default.target
