# .config/systemd/user/udiskie.service
[Unit]
Description=Udiskie
After=default.target

[Service]
ExecStart=/usr/bin/udiskie
Restart=on-failure

[Install]
WantedBy=default.target
